# files/system/etc/xdg/autostart/os-autostart.desktop
[Desktop Entry]
Type=Application
Name=OS Autostart
Exec=/usr/libexec/os-autostart
